# packaging/systemd/splitway.service
[Unit]
Description=Splitway split-DNS daemon
Documentation=https://github.com/stslex/splitway
# VPN detection needs NetworkManager (D-Bus); applying rules needs
# systemd-resolved (resolvectl). Order after both so they are present.
After=network-online.target NetworkManager.service systemd-resolved.service
Wants=network-online.target

[Service]
Type=simple
# Installed to /usr/bin by deb/rpm/pacman packages. (The Nix build wires
# its own ExecStart via the package store path; see nix/module.nix.)
# `--config` points the daemon at the persistent StateDirectory below so it
# creates its config there on first run, instead of falling back to
# /root/.config and logging a warning (mirrors nix/module.nix).
ExecStart=/usr/bin/splitway-daemon run --config /var/lib/splitway/config.json
Restart=on-failure
RestartSec=2

# systemd creates /run/splitway (0700, root) before start and removes it on
# stop. The daemon binds its 0600 control socket inside it; the 0700 dir is
# part of the socket's defense in depth (see packaging/README.md).
RuntimeDirectory=splitway
RuntimeDirectoryMode=0700

# Persistent state dir /var/lib/splitway (0700, root), created by systemd
# before start and preserved across restarts. The daemon writes/reads its
# config at <state>/config.json (see --config above); on first run it creates
# an empty config there. Mirrors `StateDirectory` in nix/module.nix.
StateDirectory=splitway
StateDirectoryMode=0700

# Unprivileged GUI/CLI access (opt-in): let a non-root user in a dedicated group
# drive the daemon without sudo. To enable, create + join the group
# (`groupadd splitway`, `usermod -aG splitway <user>`) and uncomment all four
# lines below — the socket becomes 0660 and the runtime dir 0750, both
# root:splitway. The bare `ExecStart=` first clears the default start command so
# the second line replaces it (rather than adding a second one). SECURITY: group
# membership ≈ control of system split-DNS routing (see packaging/README.md,
# "Socket security model"). Default is root-only.
#ExecStart=
#ExecStart=/usr/bin/splitway-daemon run --config /var/lib/splitway/config.json --socket-group splitway
#RuntimeDirectoryMode=0750
#SupplementaryGroups=splitway

# Runs as root deliberately. `resolvectl dns|domain|revert <iface>` performs
# privileged DNS changes that polkit otherwise gates behind interactive
# authentication — unworkable for an unattended daemon. Running as root
# sidesteps that. The privilege boundary for the unprivileged CLI is the
# Unix control socket (mode 0600); see packaging/README.md for the threat
# model. (A user service would need polkit rules granting the DNS actions.)

# On stop, systemd sends SIGTERM, which the daemon traps to revert active
# DNS rules before exiting — the system is never left half-configured.
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
